# learn/script/rc.local
#!/bin/sh
   /path/to/start.sh
   exit 0